# systemd/mlmc.timer
[Unit]
Description=Run mlm dry report weekly
[Timer]
OnCalendar=*-*-05 12:00:00
Persistent=true
[Install]
WantedBy=timers.target
